# mtb-mysql01/appserver.ini
[program:appserver]
command=/usr/bin/python /usr/local/bin/appserver.py
autostart=true
autorestart=true
user=root
