# mtb-mysql01/appserver.ini
[program:appserver]
command=/bin/bash /usr/local/bin/appserver.sh
autostart=true
autorestart=true
user=root
startretries=300
