# mtb-mysql01/mysqld.ini
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mariadb/mariadb.pid /usr/bin/mysqld_safe --basedir=/usr
autostart=true
autorestart=true
user=root
