# webserver.ini
[program:webserver]
command=/bin/bash /usr/local/bin/webserver.sh
autostart=true
autorestart=true
user=root
startretries=300
