# webserver.ini
[program:webserver]
command=/usr/bin/python /usr/local/bin/webserver.py
autostart=true
autorestart=true
user=root
startretries=300
